# layout/DEBIAN/postinst
#!/bin/sh
echo "Post install: Loading daemon..."
launchctl load /Library/LaunchDaemons/com.mikaelbo.proxyswitcherd.plist
exit 0