# layout/DEBIAN/postinst
#!/bin/sh
echo "Loading ProxySwitcher daemon..."
launchctl load /Library/LaunchDaemons/com.mbo42.proxyswitcherd.plist
exit 0